# deploy/mafex-api.service
[Unit]
Description=Mafex Room Booking API (FastAPI)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=mafex
Group=mafex
WorkingDirectory=/opt/mafex/mafexAll-main
Environment=PYTHONPATH=/opt/mafex/mafexAll-main
EnvironmentFile=/opt/mafex/mafexAll-main/.env
ExecStart=/opt/mafex/mafexAll-main/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
